# .config/systemd/user/change-gnome-theme.timer
[Unit]
Description=Change GNOME theme

[Timer]
OnCalendar=*-*-* 7:00:00
OnCalendar=*-*-* 22:00:00
Persistent=true

[Install]
WantedBy=timers.target
